# profile/airootfs/etc/systemd/system/clearwater-flathub.service
[Unit]
Description=Enable the Flathub Flatpak remote
After=NetworkManager.service
Wants=NetworkManager.service

[Service]
Type=oneshot
ExecStart=-/usr/bin/flatpak remote-add --if-not-exists --system flathub https://dl.flathub.org/repo/flathub.flatpakrepo

[Install]
WantedBy=multi-user.target
